# root/etc/services.d/nginx/run
#!/usr/bin/with-contenv sh

exec /usr/sbin/nginx -c /config/nginx.conf